# contrib/init/cspnd.conf
description "CSPN Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env CSPND_BIN="/usr/bin/cspnd"
env CSPND_USER="cspn"
env CSPND_GROUP="cspn"
env CSPND_PIDDIR="/var/run/cspnd"
# upstart can't handle variables constructed with other variables
env CSPND_PIDFILE="/var/run/cspnd/cspnd.pid"
env CSPND_CONFIGFILE="/etc/cspn/cspn.conf"
env CSPND_DATADIR="/var/lib/cspnd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # cspnd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$CSPND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run cspnd."
        echo "The setting must appear in $CSPND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"CSPN Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$CSPND_PIDDIR"
    chmod 0755 "$CSPND_PIDDIR"
    chown $CSPND_USER:$CSPND_GROUP "$CSPND_PIDDIR"
    chown $CSPND_USER:$CSPND_GROUP "$CSPND_CONFIGFILE"
    chmod 0660 "$CSPND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$CSPND_PIDFILE" \
    --chuid $CSPND_USER:$CSPND_GROUP \
    --exec "$CSPND_BIN" \
    -- \
    -pid="$CSPND_PIDFILE" \
    -conf="$CSPND_CONFIGFILE" \
    -datadir="$CSPND_DATADIR" \
    -disablewallet \
    -daemon

